# roles/cat/files/upload.service
[Unit]
Description = upload daemon

[Service]
ExecStart = /root/bin/upload.sh
Restart = always
Type = simple

[Install]
WantedBy = multi-user.target
